# awesome/overlay-livecd/usr/lib/systemd/system/livecd.service
[Unit]
Description=LiveMedia Config Script
Before=systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=/etc/manjaro/livecd

[Install]
WantedBy=multi-user.target
