# awesome/overlay-livecd/usr/lib/systemd/system/pamac.service
[Unit]
Description=Pamac
After=systemd-user-sessions.service

[Service]
Type=dbus
BusName=org.manjaro.pamac
ExecStart=/usr/bin/pamac-daemon
